# post_805/inittab
1:2345:respawn:/sbin/getty -a root --noclear 38400 tty1
#1:2345:respawn:/sbin/getty --noclear 38400 tty1